# scripts/telegram-bridge.service
[Unit]
Description=BR.AI.ON Telegram Bridge
After=network.target

[Service]
Type=simple
User=mcpgw
WorkingDirectory=/home/mcpgw/br-ai-on
ExecStart=/bin/bash /home/mcpgw/br-ai-on/scripts/telegram-bridge.sh
Restart=always
RestartSec=5
StandardOutput=append:/home/mcpgw/br-ai-on/logs/telegram-bridge.log
StandardError=append:/home/mcpgw/br-ai-on/logs/telegram-bridge.log
Environment="TMUX_TMPDIR=/tmp"

[Install]
WantedBy=multi-user.target
